# scripts/com.ainews.tracker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ainews.tracker</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/linavarbanova/Documents/claude_code_folder/ai-news-tracker/.venv/bin/python</string>
        <string>/Users/linavarbanova/Documents/claude_code_folder/ai-news-tracker/scheduler.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/linavarbanova/Documents/claude_code_folder/ai-news-tracker</string>

    <!-- Run every 3 hours (10800 seconds) -->
    <key>StartInterval</key>
    <integer>10800</integer>

    <!-- Also fire ~30 seconds after load so first run is soon -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/linavarbanova/Documents/claude_code_folder/ai-news-tracker/logs/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/linavarbanova/Documents/claude_code_folder/ai-news-tracker/logs/launchd.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
